# noten-update.timer
[Unit]
Description=Notenservice Updater

[Timer]
OnBootSec=5min
OnUnitInactiveSec=30min

[Install]
WantedBy=basic.target